# robot.service
[Unit]
Description=Robot service

[Service]
ExecStart=/home/robo/.local/bin/uv run python3 /home/robo/robocup2026-raspberrypi-kanto/main.py
Restart=always
User=robo
Group=robo
WorkingDirectory=/home/robo7


[Install]
WantedBy=multi-user.target
