# robot.service
[Unit]
Description=Robot service

[Service]
ExecStart=uv run python3 /home/robo/robocup2026-raspberrypi-program/main.py
Restart=always
User=robo
Group=robo

[Install]
WantedBy=multi-user.target
